# log-search.service
[Unit]
Description=Log Search MCP Service
After=network.target

[Service]
Type=simple
# 建议创建一个专用用户运行服务，例如 logsearch
# User=logsearch
# Group=logsearch
User=root
WorkingDirectory=/opt/log-search
# 假设二进制文件位于 /usr/local/bin，配置文件位于 WorkingDirectory
ExecStart=/usr/local/bin/log-search-mcp config.yaml
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
